# distribution/src/main/package/distributions/RedHat/src/usr/lib/systemd/system/unity-idm.service
[Unit]
Description=Unity IdM Service
After=syslog.target network.target

[Service]
Type=forking
User=unity-idm
PIDFile=/var/run/unity-idm/unity-idm-server.pid
ExecStart=/usr/sbin/unity-idm-server-start
ExecStop=/usr/sbin/unity-idm-server-stop
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
